# deploy/bot.service
[Unit]
Description=Crypto Options Bot (Delta India)
Documentation=https://github.com/yesterdays-rebel/crypto-options-bot
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=simple
User=bot
Group=bot
WorkingDirectory=/opt/crypto-options-bot
EnvironmentFile=-/opt/crypto-options-bot/.env
ExecStartPre=/usr/bin/docker compose -f /opt/crypto-options-bot/docker-compose.prod.yml pull
ExecStart=/usr/bin/docker compose -f /opt/crypto-options-bot/docker-compose.prod.yml up --remove-orphans
ExecStop=/usr/bin/docker compose -f /opt/crypto-options-bot/docker-compose.prod.yml down
Restart=on-failure
RestartSec=10
TimeoutStopSec=120
KillSignal=SIGTERM
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=read-only
PrivateTmp=true

[Install]
WantedBy=multi-user.target
